# code/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:mongod]
command=mongod --bind_ip 127.0.0.1 --port 27017 --dbpath /data/db --logpath /dev/stdout --logappend
autostart=true
autorestart=true
startretries=3
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=nobody

[program:uvicorn]
command=uv run uvicorn sifter.server:app --host 0.0.0.0 --port 8000
directory=/app/server
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=SIFTER_MONGODB_URI="mongodb://127.0.0.1:27017"
